# gas-station/standalone/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
logfile_maxbytes=10MB
logfile_backups=3

; ── PostgreSQL 16 ──────────────────────────────────────────────────────────
[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D /var/lib/postgresql/data -c listen_addresses=*
user=postgres
priority=10
autorestart=true
startretries=3
stdout_logfile=/var/log/supervisor/postgres.log
stderr_logfile=/var/log/supervisor/postgres.log

; ── Gas Station Data Generator ─────────────────────────────────────────────
[program:generator]
command=/app/venv/bin/python /app/generator/main.py
priority=20
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/var/log/supervisor/generator.log
stderr_logfile=/var/log/supervisor/generator.log
environment=POSTGRES_HOST="127.0.0.1",POSTGRES_PORT="5432",POSTGRES_DB="%(ENV_POSTGRES_DB)s",POSTGRES_USER="%(ENV_POSTGRES_USER)s",POSTGRES_PASSWORD="%(ENV_POSTGRES_PASSWORD)s"

; ── FastAPI / Uvicorn ──────────────────────────────────────────────────────
[program:api]
command=/app/venv/bin/uvicorn main:app --host 0.0.0.0 --port 8000 --app-dir /app/api
priority=20
autorestart=true
startretries=5
startsecs=3
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api.log
environment=POSTGRES_HOST="127.0.0.1",POSTGRES_PORT="5432",GAS_STATION_DB="%(ENV_POSTGRES_DB)s",GROCERY_DB="",POSTGRES_USER="%(ENV_POSTGRES_USER)s",POSTGRES_PASSWORD="%(ENV_POSTGRES_PASSWORD)s"

; ── Streamlit UI ───────────────────────────────────────────────────────────
[program:ui]
command=/app/venv/bin/streamlit run /app/ui/app.py
    --server.port=8501
    --server.address=0.0.0.0
    --server.headless=true
    --server.fileWatcherType=none
priority=30
autorestart=true
startretries=5
startsecs=5
stdout_logfile=/var/log/supervisor/ui.log
stderr_logfile=/var/log/supervisor/ui.log
environment=API_BASE_URL="http://127.0.0.1:8000"
